# app/config/docker/php-fpm-7/rc.local
#!/bin/sh -e
/usr/local/bin/fixtures.sh
exit 0